# .config/systemd/user/dunst.service
[Unit]
Description=dunst
After=graphical.target

[Service]
ExecStart=/usr/bin/dunst
Restart=always

[Install]
WantedBy=default.target
