# .config/systemd/user/wallpaper.service
[Unit]
Description=Wallpaper

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/home/markus/.fehbg

[Install]
WantedBy=window-manager.target
