# .config/systemd/user/window-manager.target
[Unit]
Description=Window Manager running
After=graphical.target
